# src/build/modules/base/files/etc/supervisor/conf.d/cron.conf
[program:cron]
command = /usr/sbin/cron -f
autostart = true
autorestart = true
stdout_events_enabled = true
stderr_events_enabled = true
