# src/build/modules/base/files/etc/supervisor/conf.d/postfix.conf
[program:postfix]
process_name = master
command = /usr/sbin/postfix -c /etc/postfix start
startsecs = 0
autostart = true
autorestart = false
stdout_events_enabled = true
stderr_events_enabled = true
